# dell-fan-control.service
[Unit]
Description=Dell fan control

[Service]
ExecStart=dell-fan-control.py
Restart=on-failure

[Install]
WantedBy=multi-user.target